# ini/cron.ini
[program:cron]
command=/usr/sbin/crond -n
