# ini/rsyslog.ini
[program:rsyslog]
command=/sbin/rsyslogd -n -c5
startsecs=5
stopwaitsecs=5
